# cloudpub/rootfs/etc/s6-overlay/s6-rc.d/cloudpub/finish
#!/command/with-contenv bashio

# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at https://mozilla.org/MPL/2.0/.

if grep -q "AuthFailed" /root/.cache/cloudpub/client.log 2>/dev/null; then
    bashio::log.error "Обнаружена ошибка: Неверный ключ API."
    echo 1 > /run/s6-linux-init-container-results/exitcode
    exec /run/s6/basedir/bin/halt
fi

if test "$1" -eq 256 ; then
  e=$((128 + $2))
else
  e="$1"
fi

if test "$e" -ne 0; then
    bashio::log.warning "Клиент CloudPub завершил работу с кодом $e (сигнал: $2)."
    echo "$e" > /run/s6-linux-init-container-results/exitcode
    exec /run/s6/basedir/bin/halt
else
    bashio::log.info "Клиент CloudPub завершил работу с кодом $e (сигнал: $2)."
    echo "$e" > /run/s6-linux-init-container-results/exitcode
    exec /run/s6/basedir/bin/halt
fi
