# cloudpub/rootfs/etc/s6-overlay/s6-rc.d/cloudpub/run
#!/command/with-contenv bashio

# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at https://mozilla.org/MPL/2.0/.

LOG_LEVEL=$(bashio::config 'log_level')

CLO="clo --conf /config/client.toml --log-level ${LOG_LEVEL}"

bashio::log.info "Запуск аддона CloudPub с уровнем логирования: ${LOG_LEVEL}..."

TOKEN=$(bashio::config 'token')
bashio::log.info "Токен загружен из конфигурации."

$CLO set token "$TOKEN"
bashio::log.info "Аккаунт CloudPub успешно привязан."

SERVICES=$(bashio::config 'services')

for ID in $(bashio::config "services|keys"); do
  TYPE=$(bashio::config "services[${ID}].type")
  LOCAL_ADDR=$(bashio::config "services[${ID}].local_addr")
  NAME=$(bashio::config "services[${ID}].name")

  bashio::log.info "Регистрация сервиса:"
  bashio::log.info "  - Тип: ${TYPE}"
  bashio::log.info "  - Локальный адрес: ${LOCAL_ADDR}"

  if bashio::var.has_value "${NAME}"; then
    bashio::log.info "  - Имя: ${NAME}"
    bashio::log.info "Выполнение команды: ${CLO} register --name ${NAME} ${TYPE} ${LOCAL_ADDR}"
    $CLO register --name "${NAME}" "${TYPE}" "${LOCAL_ADDR}"
  else
    bashio::log.info "Выполнение команды: ${CLO} register ${TYPE} ${LOCAL_ADDR}"
    $CLO register "${TYPE}" "${LOCAL_ADDR}"
  fi

  bashio::log.info "Сервис успешно зарегистрирован."
done

bashio::log.info "Запуск CloudPub в фоновом режиме..."
exec $CLO -v run
