# cloudpub/rootfs/etc/s6-overlay/s6-rc.d/clo-watchdog/run
#!/command/with-contenv bashio

# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at https://mozilla.org/MPL/2.0/.

bashio::log.info "Запуск watchdog..."

LOG_FILE="/root/.cache/cloudpub/client.log"
ERROR_PATTERN="WebSocket protocol error: Connection reset without closing handshake"

# Watchdog script to restart the CloudPub client if the data channel fails to initialize.
# Detects specific WebSocket errors and triggers a restart to recover.
stdbuf -o L tail -F "${LOG_FILE}" 2>/dev/null | while read -r line; do
  if echo "${line}" | grep -q "${ERROR_PATTERN}"; then
    bashio::log.error "Обнаружена ошибка: ${ERROR_PATTERN}. Перезапуск клиента CloudPub..."

    s6-svc -h /run/s6-rc/servicedirs/clo

    # Cleanup logs to avoid restart loops
    > "${LOG_FILE}"
  fi
done
